# puppet/modules/puppet/files/puppetd.cron
SHELL=/bin/bash
*/10 * * * *     root    ( . /etc/sysconfig/puppet ; nice -n 19 puppetd --onetime --verbose --no-daemonize --splaylimit=150 --splay $PUPPET_EXTRA_OPTS ) > /tmp/puppet_last_run.txt 2>&1
